# runit/simulator/run
#!/bin/sh

cd /usr/share/bytebeam/simulator

exec RUST_LOGS=error simulator
